# nginx/token-stats@.service
[Unit]
Description=Token Stats Dashboard Backend (port %i)
After=network.target

[Service]
Type=simple
User=skyscribe
WorkingDirectory=/home/skyscribe/srcs/token-stats/backend
Environment="RUST_LOG=info"
Environment="PORT=%i"
ExecStart=/home/skyscribe/srcs/token-stats/backend/target/release/token-stats-backend
Restart=on-failure
RestartSec=5
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
